# deploy/instagram-live-recorder.service
[Unit]
Description=Instagram Live Recorder
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/instagram-recorder/instagram-live-recorder
Environment="PATH=/home/ubuntu/instagram-recorder/venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/ubuntu/instagram-recorder/venv/bin/python main.py
Restart=always
RestartSec=10

# 로그 설정
StandardOutput=append:/home/ubuntu/instagram-recorder/logs/live-recorder.log
StandardError=append:/home/ubuntu/instagram-recorder/logs/live-recorder-error.log

# 리소스 제한
MemoryMax=1G
CPUQuota=50%

[Install]
WantedBy=multi-user.target
